# systemd/lab-api.service
[Unit]
Description=linux-ops-lab — Fake API service (HTTP :8080)
Documentation=https://github.com/devBoya/devops-up-100
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=labapp
Group=labapp
WorkingDirectory=/opt/linux-ops-lab
ExecStart=/opt/linux-ops-lab/scripts/api.sh
Restart=on-failure
RestartSec=3s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lab-api

# Light sandboxing — keep it educational, not paranoid
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
